# etc/systemd/system/flask.gunicorn.service
[Unit]
Description=Veterynary clinic application serve by Fastapi and Gunicorn 
Requires=fastapi.gunicorn.socket
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=/opt/src/fastapi-gunicorn/
ExecStart=/opt/gunicorn applicationname.wsgi
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target